# .config/autostart/xkb.desktop
[Desktop Entry]
Name=Load keyboard config
Exec=/home/tf/bin/xkb.fish
Type=Application
Comment=

