# .config/autostart/xkb.desktop
[Desktop Entry]
Name=Load keyboard config
Exec=xkb.fish
Type=Application
Comment=
Hidden=false

