# .config/autostart/xkb.desktop
[Desktop Entry]
Name=Load keyboard config
Exec=~/bin/xkb.sh
Type=Application
