# .config/autostart/xresources.desktop
[Desktop Entry]
Name=Load .Xresources
Exec=xrdb ~/.Xresources
Type=Application
